# docker/docker-php-supervisor/cron.d/root
* * * * * /usr/local/bin/php /data01/api.console.insurforce.cn/artisan schedule:run
